# bridge/systemd/milo-bridge.service
[Unit]
Description=Milo bridge - robot drivers, gait engine, knowledge graph, brain link
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dama
ExecStart=/home/dama/.venvs/milo/bin/python -m milo_bridge.main
Restart=always
RestartSec=3
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
